# roles/01-system/templates/etc/cron.d/sync-skel-userdir
0 0 * * *  root bash /opt/sync-skel-userdir.sh
