# roles/01-system/templates/etc/cron.d/update-stapling
15 0 * * * root bash /opt/update-stapling.sh > /dev/null 2>&1
