# roles/01-system/templates/etc/cron.d/update-dh
15 0 * * * root bash /opt/update-dh.sh > /dev/null 2>&1
