# roles/01-system/templates/etc/cron.d/update-ssl
0 0 1 * * root bash /opt/update-ssl.sh > /dev/null
